# systemd/system/mctp-spi-ctrl.service
[Unit]
Description=MCTP SPI control daemon
After=mctp-spi-demux.service
Requires=mctp-spi-demux.service

[Service]
Restart=always
Environment=MCTP_SPI_CTRL_OPTS=null
EnvironmentFile=-/etc/default/mctp
ExecStart=/usr/bin/mctp-spi-ctrl $MCTP_SPI_CTRL_OPTS
SyslogIdentifier=mctp-spi-ctrl

[Install]
WantedBy=multi-user.target
